# kthresher.cron.daily
#!/bin/sh
# Regular cron job for the mykthresher package
#

[ -x /usr/sbin/kthresher ] && /usr/sbin/kthresher
