# debian/kthresher.cron.daily
#!/bin/sh

[ -x /usr/bin/kthresher ] && /usr/bin/kthresher
